# rootfs/etc/services.d/miner/run
#!/bin/sh

echo "* generating block..."

dash-cli -conf=/etc/conf/node-1/dash.conf generate 1
sleep 10